# archlinux/config-graphic/config-graphic.install
post_install() {
    # Setup screen lock on suspend
    systemctl enable suspend.service
}
